# src/root/etc/s6-overlay/s6-rc.d/sshd/run
#!/command/execlineb -P
# フォアグラウンドで sshd を exec
importas DEBUG DEBUG

define -s SSHD_RUN     "/usr/sbin/sshd -D -e"
define -s SSHD_RUN_D   "/usr/sbin/sshd -D -e -d"

foreground { ssh-keygen -A }

# # Redirect stderr to stdout.
fdmove -c 2 1

$SSHD_RUN
# $SSHD_RUN_D